# image/stage5-nexus/01-nexus-services/files/nexus-gateway.service
[Unit]
Description=NEXUS Gateway Daemon
After=network-online.target nexus-node-agent.service
Wants=network-online.target

[Service]
Type=simple
User=nexus
Group=nexus
WorkingDirectory=/opt/nexus
ExecStart=/usr/bin/python3 /opt/nexus/agents/gateway.py
Restart=on-failure
RestartSec=10
StandardOutput=append:/opt/nexus/logs/gateway.log
StandardError=append:/opt/nexus/logs/gateway.log

[Install]
WantedBy=multi-user.target
