# image/stage5-nexus/01-nexus-services/files/nexus-codebase/deploy/node-agent.service
[Unit]
Description=NEXUS Node Agent
After=network.target

[Service]
User=mhuraibi
WorkingDirectory=/opt/nexus/agents
EnvironmentFile=/opt/nexus/agents/.env
ExecStart=/usr/bin/python3 /opt/nexus/agents/node_agent.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
